# dot_config/systemd/user/lockscreen-refresh.service.tmpl
[Unit]
Description=Render lockscreen composite + mirror to greeter (pycairo + Pango)

[Service]
Type=oneshot
Environment="PATH={{ .chezmoi.homeDir }}/.local/share/mise/shims:/usr/local/bin:/usr/bin:/bin"
Environment="DISPLAY=:0"
ExecStart={{ .chezmoi.homeDir }}/.local/share/mise/shims/python3 %h/.config/i3/scripts/lock-render.py --lock-image %h/.cache/lockscreen/_compose.jpg
StandardOutput=journal
StandardError=journal
